# services/incus-templates/files/systemd/personal-server-incus-template-refresh.timer
[Unit]
Description=personal-server Incus template weekly refresh (Mon 03:00)

[Timer]
# Aligned with the rest of the weekly update wave: 03:00 = update,
# 04:00 = reboot.
OnCalendar=Mon *-*-* 03:00:00
Persistent=true
AccuracySec=1min
Unit=personal-server-incus-template-refresh.service

[Install]
WantedBy=timers.target
