# shared/systemd/bootc-fetch-apply-updates.timer.d/10-personal-server.conf
# Dropin: pin bootc upgrade to the weekly Monday 03:00 update wave.
#
# Default upstream values fire every 8h with a 2h random delay (OnBootSec=1h,
# OnUnitInactiveSec=8h, RandomizedDelaySec=2h). Empty assignments before
# OnCalendar are required so systemd clears the inherited list before
# applying ours — otherwise the new rule is appended to the existing ones.
[Timer]
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=
OnCalendar=Mon *-*-* 03:00:00
Persistent=true
