# shared/systemd/personal-server-reboot.timer
[Unit]
Description=personal-server scheduled full reboot (Monday 04:00)

[Timer]
# Sequenced after the 03:00 update wave (bootc + watchtower + incus
# templates) so the reboot picks up everything that landed during the hour.
OnCalendar=Mon *-*-* 04:00:00
Persistent=true
AccuracySec=1min
Unit=personal-server-reboot.service

[Install]
WantedBy=timers.target
